# emacs/.config/systemd/user/emacs.service.d/override.conf
[Service]
Environment=GDK_BACKEND=wayland
