# ssh/.config/systemd/user/ssh-add-keys.service
# Load this machine's YubiKey-backed sk keys into the ssh-agent once per login.
#
# Loading an sk key only registers its on-disk credential handle — it never
# contacts the YubiKey (touch/PIN happen at *use* time), so this runs headless
# and even with the key unplugged. Keys must be in the local agent for SSH
# agent forwarding (e.g. to am-jallen / chonky) to broker them.
#
# Enable once per machine:  systemctl --user enable --now ssh-add-keys.service
[Unit]
Description=Load YubiKey sk keys into the ssh-agent
After=ssh-agent.socket
Requires=ssh-agent.socket
ConditionPathExists=%h/.ssh/id_ed25519_sk_gitsign

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=SSH_AUTH_SOCK=%t/ssh-agent.socket
# Load only the sk keys that exist on this machine — not every machine owns all
# three, and ssh-add fails the whole unit if any listed file is missing.
ExecStart=/bin/sh -c 'k=$(ls %h/.ssh/id_ed25519_github_sk %h/.ssh/id_ed25519_sk_gitsign %h/.ssh/id_ed25519_sk 2>/dev/null); [ -n "$k" ] && exec ssh-add $k'

[Install]
WantedBy=default.target
