# ssh/.config/systemd/user/ssh-add-keys.service
# Load this machine's YubiKey-backed sk keys into the ssh-agent once per login.
#
# Loading an sk key only registers its on-disk credential handle — it never
# contacts the YubiKey (touch/PIN happen at *use* time), so this runs headless
# and even with the key unplugged. Keys must be in the local agent for SSH
# agent forwarding (e.g. to am-jallen / chonky) to broker them.
#
# The set of sk keys varies per host, so the list lives in a host-local file
# ~/.ssh/agent-keys (one key per line, relative to ~/.ssh; '#' comments and
# blank lines ignored; missing files skipped). That file is not stowed, so it
# stays out of the dotfiles repo. The unit no-ops on hosts without it.
#
# Enable once per machine:  systemctl --user enable --now ssh-add-keys.service
[Unit]
Description=Load YubiKey sk keys into the ssh-agent
After=ssh-agent.socket
Requires=ssh-agent.socket
ConditionPathExists=%h/.ssh/agent-keys

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=SSH_AUTH_SOCK=%t/ssh-agent.socket
# Logic lives in a script — inline shell here would collide with systemd's own
# ${...} / %% expansion in ExecStart.
ExecStart=%h/.local/bin/load-ssh-agent-keys

[Install]
WantedBy=default.target
